# server-backend/systemd/dwargonmc-backend.service
[Unit]
Description=DwargonMC Launcher Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/dwargonmc-launcher/server-backend
Environment=PORT=2121
Environment=BIND_HOST=127.0.0.1
Environment=PUBLIC_URL=https://dwargonmc-sync.petershub.xyz
Environment=MAP_TARGET=http://127.0.0.1:8888
# Optional when MAP_TARGET is protected by Cloudflare Access:
# Environment=MAP_ACCESS_CLIENT_ID=
# Environment=MAP_ACCESS_CLIENT_SECRET=
# Or arbitrary JSON headers:
# Environment=MAP_REQUEST_HEADERS={"x-map-token":"secret"}
Environment=MC_HOST=127.0.0.1
Environment=MC_PORT=25565
ExecStart=/usr/bin/node /opt/dwargonmc-launcher/server-backend/server.js
Restart=always
RestartSec=5
User=dwargonmc
Group=dwargonmc

[Install]
WantedBy=multi-user.target
